# scripts/com.weaper.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS Launch Agent for Weaper Sync Server.

  BEFORE INSTALLING — update the two placeholder values below:
    1. ProgramArguments[1]: replace "/path/to/Weaper/local-server/src/index.js"
       with the absolute path to the index.js file on your Mac.
       Example: /Users/yourname/Projects/Weaper/local-server/src/index.js

    2. EnvironmentVariables > UPLOAD_DIR: replace the default path with the
       directory where REAPER reads its media files on your Mac.
       Example: /Users/yourname/Documents/REAPER Media/Weaper

  Install with:
    cp scripts/com.weaper.server.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.weaper.server.plist

  Uninstall with:
    launchctl unload ~/Library/LaunchAgents/com.weaper.server.plist
    rm ~/Library/LaunchAgents/com.weaper.server.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.weaper.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
  <!-- UPDATE: replace with your absolute path to local-server/src/index.js -->
    <string>/path/to/Weaper/local-server/src/index.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>3000</string>
    <key>UPLOAD_DIR</key>
    <!-- UPDATE: replace with your REAPER media directory, e.g. /Users/yourname/Documents/REAPER Media/Weaper -->
    <string>/Users/user/Documents/REAPER Media/Weaper</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/weaper-server.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/weaper-server-error.log</string>
</dict>
</plist>
